# sxhkd/.config/autostart/sxhkd.desktop
[Desktop Entry]
Type=Application
Exec=bash -c 'sxhkd'
Hidden=false
Autostart-enabled=true
Name=keybind
Comment=hotkey daemon
