# sxhkd/.config/autostart/sxhkd.desktop
[Desktop Entry]
Type=Application
Exec=bash -c 'sxhkd & disown'
Hidden=false
X-GNOME-Autostart-enabled=true
Name=keybind
Comment=hotkey daemon
