# sxhkd/.config/autostart/sxhkd.desktop
[Desktop Entry]
Type=Application
Exec=sxhkd
Hidden=false
Terminal=false
X-GNOME-Autostart-enabled=true
Name=keybind
Comment=hotkey daemon
